# media-video/casparcg-server/files/casparcg-server.service
[Unit]
Description=CasparCG Server Service

[Service]
Type=exec
ExecStart=/usr/bin/casparcg %h/.casparcg/casparcg.config
